# tutu.service
[Unit]
Description=Foo

[Service]
ExecStart=/usr/sbin/luxeriadoorbellclient

[Install]
WantedBy=multi-user.target
